# files/etc/init.d/ip-qos
#!/bin/sh /etc/rc.common
START=99
start() {
    /usr/bin/ip-qos-apply
}
